# crontab
0 18 * * * cd /app && npm run daily-sync >> /app/logs/afl-sync.log 2>&1