# crontab
0 4 * * * cd /app && npm run daily-sync >> /var/log/afl-sync.log 2>&1